# crontab
*/15 * * * * cd ~/.openclaw/workspace/crypto-edge && node services/trading-engine.js >> ~/.openclaw/workspace/crypto-edge/logs/cron.log 2>&1
